# debian/uplinkmgr.service
[Unit]
Description=uplinkmgr multi-WAN uplink monitor daemon
After=network.target dhcpcd.service

[Service]
Type=simple
ExecStart=/usr/sbin/uplinkmgr
ExecStop=/bin/kill -TERM $MAINPID
Restart=on-failure
RestartSec=10s
RuntimeDirectory=uplinkmgr
RuntimeDirectoryMode=0755

[Install]
WantedBy=multi-user.target
